# apps/planner/paper-device/paperos-cache-refresh.service
[Unit]
Description=PaperOS Planner cache refresh

[Service]
Type=oneshot
WorkingDirectory=/home/root/paperos
ExecStart=/home/root/paperos/refresh-cache.sh
