# apps/planner/paper-device/paperos-cache-refresh.timer
[Unit]
Description=Refresh the PaperOS Planner cache every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=paperos-cache-refresh.service

[Install]
WantedBy=timers.target
